# arcanasnapshotbot-vampire/arcana-sync.timer
[Unit]
Description=Daily Auto-Sync for Arcana Sovereign Node

[Timer]
# Run 5 minutes after boot
OnBootSec=5min
# Run every 24 hours
OnUnitActiveSec=24h
Unit=arcana-sync.service

[Install]
WantedBy=timers.target
